# configs/headset/headset-volume.service
[Unit]
Description=Set headset volume to 70

[Service]
Type=oneshot
ExecStart=/usr/bin/headsetcontrol -s 70
